# edexc/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:systemd_init]
command=/user/sbin/init

[program:edex_start]
command=edex start
