# edexc/etc/systemd/edex_start.service
[Unit]
Description=Start edex server

[Service]
Type=oneshot
ExecStart=/scripts/edex_start.sh

[Install]
WantedBy=multi-user.target
