# edexc/etc/systemd/edex_start.service
[Unit]
Description=Start edex server

[Service]
Type=oneshot
ExecStart=/usr/bin/edex start

[Install]
WantedBy=multi-user.target
